# service/chisel/run
#!/bin/sh
sleep 90
chisel server --auth password --port 17860
